# printer.service
[Unit]
Description=Printer Service

[Service]
ExecStart=/tmp/CST/CSTFP
Restart=on-failure

[Install]
WantedBy=multi-user.target
